# ansible/files/quadlet/portainer.container.j2
[Unit]
Description=Portainer - Container management UI for PrivateBox
Documentation=https://docs.portainer.io/
Wants=network-online.target podman.socket
After=network-online.target podman.socket
{% if use_system_quadlet | default(true) %}
Before=multi-user.target
{% endif %}

[Container]
# Container image
Image={{ container_image_registry | default('docker.io') }}/{{ portainer_image }}:{{ portainer_version }}
ContainerName=portainer

# Network configuration
PublishPort={{ ansible_default_ipv4.address }}:{{ portainer_web_port }}:9443

# Volume mounts
Volume=/run/podman/podman.sock:/var/run/docker.sock:z
Volume={{ portainer_data_dir }}:/data:{{ volume_mount_options | default('Z') }}
Volume=snippets:/snippets:{{ volume_mount_options | default('Z') }}
Volume=/etc/privatebox/certs:/certs:ro,z

# Environment variables
Environment="TZ={{ timezone | default('UTC') }}"

# Container arguments for SSL
Exec=--ssl --sslcert /certs/privatebox.crt --sslkey /certs/privatebox.key

# Security settings
SecurityLabelDisable={{ container_security_label_disable | default('false') | lower }}
NoNewPrivileges={{ container_no_new_privileges | default('true') | lower }}
{% if container_drop_capabilities is defined and container_drop_capabilities | length > 0 %}
{% for cap in container_drop_capabilities %}
DropCapability={{ cap }}
{% endfor %}
{% endif %}

# Resource limits
{% if portainer_memory_limit is defined %}
MemoryLimit={{ portainer_memory_limit }}
{% endif %}
{% if portainer_cpu_quota is defined %}
CPUQuota={{ portainer_cpu_quota }}
{% endif %}

# Health check
HealthCmd=wget --spider --no-check-certificate -q -O /dev/null https://127.0.0.1:9443/api/status
HealthInterval={{ health_check_interval | default('30s') }}
HealthRetries={{ health_check_retries | default(3) }}
HealthStartPeriod={{ health_check_start_period | default('60s') }}
HealthTimeout={{ health_check_timeout | default('10s') }}

# Logging
{% if container_log_driver is defined %}
LogDriver={{ container_log_driver }}
{% endif %}

# Pull policy
Pull={{ container_pull_policy | default('missing') }}

# Auto-update disabled - manual updates only via Semaphore

[Service]
# Restart policy
Restart={{ container_restart_policy | default('always') }}
RestartSec={{ container_restart_sec | default(30) }}

# Systemd service settings
TimeoutStartSec={{ systemd_timeout_start_sec | default(300) }}
TimeoutStopSec={{ systemd_timeout_stop_sec | default(90) }}

# Run pre-start commands if needed
{% if portainer_data_dir is defined %}
ExecStartPre=/bin/mkdir -p {{ portainer_data_dir }}
{% endif %}

[Install]
# Enable service at boot
WantedBy=multi-user.target default.target
